# django/configuration_files/uwsgi.conf
description "uWSGI instance to serve mysite"

start on runlevel [2345]
stop on runlevel [!2345]

setuid demo
setgid www-data

exec /usr/local/bin/uwsgi --emperor /etc/uwsgi/sites
